# bz-server/templates/default/load-user-jobs.conf.erb
start on (local-filesystems and net-device-up)

script
cat /etc/passwd | while read line
do
    user=`echo $line | cut -d: -f1`
    home=`echo $line | cut -d: -f6`
    if [ -d "$home/.init" ]
        then
        su $user -c "start $user"
    fi
done
end script
