# bz-server/files/default/master.conf.erb
start on (started user-jobs)

stop on (stopping user-jobs)
